# apps/hermes-suite/app/programs/dashboard.conf
[program:hermes-dashboard]
command=/opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0 --port 9119 --insecure --no-open
directory=/opt/hermes
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
