# apps/hermes-suite/app/programs/webui.conf
[program:hermes-webui]
command=/opt/hermes/.venv/bin/python /opt/hermes-webui/server.py
directory=/opt/hermes-webui
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
